# deploy/network/backcountry-broadcast-network.service
[Unit]
Description=Backcountry Broadcast Wi-Fi fallback
Wants=NetworkManager.service
After=NetworkManager.service
Before=backcountry-broadcast.service

[Service]
Type=oneshot
Environment=NOMADSCREEN_STORAGE_ROOT=/srv/backcountry-broadcast
ExecStart=/opt/backcountry-broadcast/deploy/network/backcountry-broadcast-network.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
