# deploy/nomadscreen-filebrowser.service
[Unit]
Description=Backcountry Broadcast File Browser
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/srv/nomadscreen/filebrowser
ExecStart=/usr/local/bin/filebrowser --address 0.0.0.0 --port 8081 --root /home/pi/media --database /srv/nomadscreen/filebrowser/filebrowser.db
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
